# etc/service/squid/run
#!/bin/bash

exec /usr/sbin/squid -YC -f /etc/squid/squid.conf --foreground
